# conf/crontab
0 * * * * /usr/bin/php /var/www/html/jupiter/rss2db/rss2db.php >/tmp/crontab.log 2>&1
