# conf/crontab
0 * * * * root /usr/bin/php /var/www/rss2db/rss2db.php
